# bin/ze-filter.service
[Unit]
Description=ze-filter mail filtering service
#After=syslog.target network.target
#Conflicts=postfix.service exim.service
#Wants=sm-client.service
#StartLimitIntervalSec=0

[Service]
Type=forking
PIDFile=/var/run/ze-filter/ze-filter.pid
#Environment=SENDMAIL_OPTS=-q1h
EnvironmentFile=-/etc/sysconfig/ze-filter
ExecStartPre=-/var/ze-filter/cdb/get-urlbl
ExecStartPre=-/var/ze-filter/cdb/get-bayes
ExecStart=/usr/sbin/ze-filter

[Install]
#WantedBy=multi-user.target
#Also=sm-client.service

