# bin/ze-greyd.init.in
#!/bin/sh
#
# @PACKAGE@ @VERSION@
#
# chkconfig: 2345 80 30
# description: ze-greyd - Joe's ze-greyd
#		filtre de messagerie

### BEGIN INIT INFO
# Provides:                     ze-greyd
# Required-Start:               $remote_fs
# Required-Stop:
# X-UnitedLinux-Should-Start:
# X-UnitedLinux-Should-Stop:
# Default-Start:                3 5
# Default-Stop:
# Description:                  Start ze-filter (filtre mail)
### END INIT INFO

#
# VARIABLES
#
ZEGREYD_HOME=@prefix@/sbin

ZEGREYD_BIN=ze-greyd

OPTIONS=" "

WORKROOT=@WORKROOT@
STATEDIR=@STATEDIR@
CONFDIR=@CONFDIR@


# Get config.

test -f /etc/sysconfig/ze-greyd && . /etc/sysconfig/ze-greyd
test -f /etc/default/ze-greyd   && . /etc/default/ze-greyd


GREYD_PIDFILE=$STATEDIR/ze-greyd.pid
SOCKFILE=$STATEDIR/ze-greyd.sock
CFFILE=$CONFDIR/ze-filter.cf

#
# Get unix socket file from configuration file
#
#TFILE=`grep "^SOCKET" $CFFILE | awk '{print $2}'`
#if [ -n "$TFILE" ];
#then
#  TFILE=`echo $TFILE | sed s/local://`
#  TFILE=`echo $TFILE | sed s/unix://`
#fi
#
#if [ -n "$TFILE" ] ;
#then
#  SOCKFILE=$TFILE
#fi

#
# Get PIDFILE file from configuration file
#
#TFILE=`grep "^GREYD_PID_FILE" $CFFILE | awk '{print $2}'`
#if [ -n "$TFILE" ] ;
#then
#  GREYD_PIDFILE=$TFILE
#fi

#
#
#
USR_ID=@USER@
GRP_ID=@GROUP@

#
# Initializing process runtime state files
#
#TPID=`grep "^GREYD_PID_FILE" $CFFILE | awk '{print $2}'`
#if [ -n "$TPID" ] ;
#then
#  # echo "Setting GREYD_PIDFILE to $TPID"
#  GREYD_PIDFILE=$TPID
#fi

#TSOCK=`grep "^UNIX_SOCK" $CFFILE | awk '{print $2}'`
#if [ -n "$TSOCK" ] ;
#then
#  # echo "Setting UNIX_SOCK to $TSOCK"
#  SOCKFILE=$TSOCK
#fi

#
# MACROS
#

check_and_kill()
{
  PID=`head -1 $1`
  kill -0 $PID > /dev/null 2>&1
  [ $? -eq 0 ] && kill $PID
}

setupdirs ()
{
  if [ -n "$WORKROOT" -a ! -d $WORKROOT ] ; then
    mkdir -p $WORKROOT
    chmod 755 $WORKROOT
    if [ -n "$USR_ID" ] ; then
      chown $USR_ID $WORKROOT
    fi
    if [ -n "$GRP_ID" ] ; then
      chgrp $GRP_ID $WORKROOT
    fi
  fi

  if [ -n "$STATEDIR" -a ! -d $STATEDIR ] ; then
    mkdir -p $STATEDIR
    chmod 755 $STATEDIR
    if [ -n "$USR_ID" ] ; then
      chown $USR_ID $STATEDIR
    fi
    if [ -n "$GRP_ID" ] ; then
      chgrp $GRP_ID $STATEDIR
    fi
  fi
}

cleanpidfile ()
{
  run=`@PSBIN@ @PSOPT@ | grep $ZEGREYD_BIN  | grep -c -v grep`
  if [ "$run" = "0" ] ; then
    if [ -f $GREYD_PIDFILE ] ; then
      echo "Removing old PID file..."
      rm -f $GREYD_PIDFILE
    fi
  fi
}

cleansockfile ()
{
  run=`@PSBIN@ @PSOPT@ | grep $ZEGREYD_BIN  | grep -c -v grep`
  if [ "$run" = "0" ] ; then
    if [ -r $SOCKFILE ] ; then
      echo "Removing old SOCK file..."
      rm -f $SOCKFILE
    fi
  fi
}

removesockfile ()
{
  if [ -r $SOCKFILE ] ; then
    echo "Removing old SOCK file..."
    rm -f $SOCKFILE
  fi
}

start ()
{
  cleanpidfile
  cleansockfile

  if [ -f $ZEGREYD_HOME/$ZEGREYD_BIN ] ; then
    echo "Starting $ZEGREYD_BIN...."
    $ZEGREYD_HOME/$ZEGREYD_BIN $OPTIONS 
    echo ".... $ZEGREYD_BIN started !"
  fi
}

stop () {
  echo "Stopping $ZEGREYD_BIN...."
  i=10
  while [ $i -gt 0 -a -f $GREYD_PIDFILE ]; do
    kill `head -1 $GREYD_PIDFILE`
    i=`expr $i - 1`
    sleep 1
  done
  cleanpidfile
  cleansockfile
  echo " OK !"
}

restart () 
{
  stop
  sleep 2
  start
}

#
# Real script begins here
#
setupdirs

case "$1" in
'start')
	start
	;;

'stop')
	stop
	;;

'restart')
	restart
	;;

*)
	echo "Usage: $0 { start | stop | restart | help }"
        cat <<EOF
start      - start ze-greyd
stop       - stop ze-greyd
restart    - stop and start ze-greyd
             (send SIGUSR2 to the daemon)
help       - this screen

EOF
	exit 1
	;;
esac

exit 0

