# config/startup/ze-greyd.service.in
[Unit]
Description=ze-filter mail filtering service

[Service]
Type=forking
PIDFile=/var/run/ze-filter/ze-greyd.pid
Environment=OPTS=
EnvironmentFile=-/etc/sysconfig/ze-greyd
ExecStart=/usr/sbin/ze-greyd $OPTIONS
ExecStop=/usr/bin/kill $PIDFile 

[Install]
WantedBy=multi-user.target

